# recipes-bsp/redeem/redeem/redeem.service
[Unit]
Description=The Replicape Dameon

[Service]
Type=oneshot
ExecStart=/usr/bin/redeem.py

[Install]
WantedBy=basic.target

