# recipes-bsp/redeem/redeem/redeem.service
[Unit]
Description=The Replicape Dameon

[Service]
Type=oneshot
ExecStart=/usr/src/redeem/software/Redeem.py

[Install]
WantedBy=basic.target

