# recipes-bsp/mjpg-streamer/files/mjpg-streamer.path
[Unit]
Description=trigger for webcam path present

[Path]
PathExists=/dev/video0

[Install]
WantedBy=multi-user.target

